# udev/10-perfmon.rules
# rules for performance monitor
# robcranfill

SUBSYSTEMS=="usb" RUN+="/home/rob/proj/practiceMonitor/udev/trigger.sh"
