# udev/10-pracmon.rules
# udev rules for Practice Monitor
# robcranfill

SUBSYSTEMS=="usb" ACTION=="change" RUN+="/home/rob/proj/practiceMonitor/udev/trigger.sh"
